# net-p2p/qbittorrent-nox/files/qbittorrent-nox-r1.initd
#!/sbin/openrc-run

command="${QBITTORRENT_NOX_BINARY}"
pidfile="${QBITTORRENT_NOX_PIDFILE}"
command_args="${SSHD_OPTS} -o PidFile=${pidfile} -f ${SSHD_CONFIG}"
